# config/cronjobs/certbot.cron
1 */12 * * *    sh /certbot.entrypoint.sh